# crawler.service.in
[Unit]
Description=AMO Crawler 

[Service]
Type=simple
ExecStart=/root/amo_crawler_run.sh @node_addr@ @chain_id@
WorkingDirectory=/root/amo-crawler
SyslogFacility=daemon
SyslogLevel=info
Restart=on-failure
LimitNOFILE=500000

[Install]
WantedBy=multi-user.target

