# server.service.in
[Unit]
Description=AMO API Server 

[Service]
Type=simple
ExecStart=/root/amo-api-server/bin/server
WorkingDirectory=/root/amo-api-server
SyslogFacility=daemon
SyslogLevel=info
Restart=on-failure
LimitNOFILE=500000

[Install]
WantedBy=multi-user.target

